# deploy/systemd/mcp-agent-mail.service
[Unit]
Description=MCP Agent Mail Service
After=network.target

[Service]
Type=simple
User=mcp-agent
Group=mcp-agent
# Path to binary
ExecStart=/usr/local/bin/mcp-agent-mail serve
# Directory where 'data' folder is located/mounted relative to CWD if implied,
# OR we should set working directory.
WorkingDirectory=/var/lib/mcp-agent-mail
# Environment variables
EnvironmentFile=/etc/default/mcp-agent-mail
# Restart on failure
Restart=on-failure
RestartSec=5
# Limits
LimitNOFILE=65536

[Install]
WantedBy=multi-user.target
