# deploy/systemd/mcp-agent-mail.service
[Unit]
Description=MCP Agent Mail Unified Service
After=network.target

[Service]
Type=simple
User=mcp
Group=mcp
# Adjust path to your binary
ExecStart=/opt/mouchak/bin/mcp-agent-mail serve http
# Adjust path to your env file
EnvironmentFile=/etc/mouchak/env
Restart=always
RestartSec=5
StandardOutput=journal
StandardError=journal
SyslogIdentifier=mcp-agent-mail

[Install]
WantedBy=multi-user.target
